# zfs_scrub/zfs-scrub@.service
[Unit]
Description=zpool scrub on %i
Requires=zfs.target
After=zfs.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/zpool scrub %i
